# data/openrc/turnstiled
#!/sbin/openrc-run

pidfile="/run/turnstiled.pid"
command="/usr/bin/turnstiled"
command_background=true
